# scripts/com.stockscouthub.edinet-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stockscouthub.edinet-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/tomomi/Documents/stock-scout-hub/scripts/publish-edinet-snapshot.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tomomi/Documents/stock-scout-hub</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>45</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/tomomi/Documents/stock-scout-hub/logs/edinet-refresh.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tomomi/Documents/stock-scout-hub/logs/edinet-refresh.launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
